# service/lighttpd/run
#!/bin/sh
/usr/sbin/lighttpd -Df /etc/lighttpd/lighttpd.conf
